# packages/opb/opb-32-py3-none-any.whl/opb-32.data/data/opb/opb.service
[Unit]
Description=object programming daemon
After=multi-user.target

[Service]
Type=forking
KillMode=mixed
User=opb
Group=opb
DynamicUser=yes
StateDirectory=opb
LogsDirectory=opb
CacheDirectory=opb
ExecStart=/usr/local/bin/opbd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
